# linux/etc/systemd/system/makima.service.d/override.conf
[Service]
User=mark
Group=input
Environment=MAKIMA_CONFIG=/home/mark/.config/makima
