# voxtype/.config/systemd/user/voxtype.service
[Unit]
Description=Voxtype push-to-talk voice-to-text daemon
Documentation=https://voxtype.io
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
# Use the ONNX build so the Parakeet engine configured in ~/.config/voxtype/config.toml works.
ExecStart=/usr/lib/voxtype/voxtype-onnx-avx2 daemon
Restart=on-failure
RestartSec=5

# Ensure we have access to the display
Environment=XDG_RUNTIME_DIR=%t

[Install]
WantedBy=graphical-session.target
